# docker-images/slurmbase/etc/supervisord.d/glusterd.conf
[program:glusterd]
command=/usr/sbin/glusterd --no-daemon
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
